# examples/config/config-exclude.list
/etc/ld.so.cache
/etc/adjtime
/etc/config-archive
/etc/lvm/cache/
/etc/lvm/backup/
/etc/lvm/archive/
/etc/ssh/ssh_host_*
/etc/blkid.tab
/etc/blkid.tab.old
/etc/mtab
/etc/passwd-
/etc/group-
/etc/shadow-
/etc/gshadow-
/etc/runlevels/install
/etc/local.d/install.start
/etc/runlevels/default/cf-serverd
/etc/runlevels/default/cf-execd
/etc/runlevels/default/cf-monitord
/etc/paludis/keywords.conf.d/cfengine.conf
/etc/paludis/use.conf.d/cfengine.conf
/etc/init.d/cf-execd
/etc/init.d/cf-monitord
/etc/init.d/cf-serverd
/etc/portage/savedconfig
/etc/ca-certificates.conf
/etc/conf.d/net.eth0
/etc/cf-template
/etc/ssl/postfix/server.*
/etc/mail/aliases.db
/etc/resolv.conf
